# systemd/easun-poller.service
[Unit]
Description=Easun Modbus Poller
After=network.target

[Service]
Type=simple
WorkingDirectory=/opt/eamoon
EnvironmentFile=/opt/eamoon/.env
ExecStart=/opt/eamoon/.venv/bin/python /opt/eamoon/modbus_service.py
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
